# .docker/php/supervisord.conf
[supervisord]
logfile = /dev/null
# loglevel = error
log_file_backups = 7
logfile_maxbytes = 100000000
nodaemon = true

[unix_http_server]
file = /tmp/supervisor.sock
chmod = 0700
username = laravel
password = laravel

[supervisorctl]
serverurl = unix:///tmp/supervisor.sock
username = laravel
password = laravel

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:php]
command = /usr/local/sbin/php-fpm
stdout_logfile = /var/www/html/storage/logs/php-fpm.log
stderror_logfile = /var/www/html/storage/logs/php-fpm.log
priority = 10
autorestart = true
autostart = true

[program:workers]
command = php artisan queue:work --tries=3 --sleep=3 --timeout=120 --max-time=3600
process_name=%(program_name)s_%(process_num)02d
stdout_logfile = /var/www/html/storage/logs/workers.log
stderror_logfile = /var/www/html/storage/logs/workers.log
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=2
redirect_stderr=true
stopwaitsecs=3600
